# crontab
0 */1 * * * server /opt/python/sähköchecker.py